# files/apache2.conf
[program:apache2]
command=apachectl -DFOREGROUND
autostart=true
autorestart=true
